# server/goweb.service
[Unit]
Description=goweb

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/genomdb/build/main

[Install]
WantedBy=multi-user.target
